# supervisor/validators/validator3.conf
[program:coq-devnet-validator]
command=/usr/local/bin/coq --datadir /root/coq-devnet/node  --syncmode full --cache 8192 --port 33710  --networkid 16390   --allow-insecure-unlock  --http --http.addr "0.0.0.0" --http.port "27320" --http.vhosts "*" --http.corsdomain "*" --ws --ws.addr "0.0.0.0" --ws.port "27330" --ws.origins "*" --http.api "eth,net,txpool,web3" --ws.api "eth,net,txpool,web3" --rpc.gascap 50000000  --maxpeers 50 --rpc.txfeecap 0 -unlock '803a106aBA7893e9785b9FDAF6CB1A9e7950B76E' --password /root/coq-devnet/password --mine  --miner.gaslimit 30000000 --nodiscover  --txlookuplimit 0 --config /root/coq-devnet/config.toml  --nat extip:192.168.9.4
stopsignal=INT
stopwaitsecs=500
stderr_logfile=/var/log/supervisor/coq-devnet.err.log
stdout_logfile=/var/log/supervisor/coq-devnet.out.log